# supervisor-cron.conf
[program:cronjob]
command = cron -f
autostart=true
autorestart=true
stdout_logfile=/var/log/cronjob-supervisor.log
stderr_logfile=/var/log/cronjob-supervisor.log
